# infrastructure/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgres]
command=/usr/lib/postgresql/14/bin/postgres -D /var/lib/postgresql/14/main -c config_file=/etc/postgresql/14/main/postgresql.conf
user=postgres
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres.err

[program:redis]
command=redis-server
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.err

[program:minio]
command=minio server /data --console-address ":9001"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/minio.log
stderr_logfile=/var/log/supervisor/minio.err

[program:backend]
command=node dist/main
directory=/app/backend
autostart=true
autorestart=true
environment=NODE_ENV="production",DB_HOST="localhost",DB_PORT="5432",REDIS_HOST="localhost",JWT_SECRET="edumap_super_secret_jwt_key_2026"
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend.err

[program:ai-service]
command=uvicorn main:app --host 0.0.0.0 --port 8000
directory=/app/ai-service
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/ai.log
stderr_logfile=/var/log/supervisor/ai.err

[program:frontend]
command=npm run start
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/frontend.log
stderr_logfile=/var/log/supervisor/frontend.err

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err
