# infrastructure/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

# NOTE: PostgreSQL and Redis are managed by system services (started in hf_entrypoint.sh).
# They are NOT managed by supervisord to avoid permission and daemonize conflicts.

[program:minio]
command=minio server /data --console-address ":9001"
priority=100
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/minio.log
stderr_logfile=/var/log/supervisor/minio.err
environment=MINIO_ROOT_USER="admin",MINIO_ROOT_PASSWORD="password123"

[program:backend]
command=bash /app/scripts/start_backend.sh
directory=/app/backend
priority=200
autostart=true
autorestart=true
startretries=50
startsecs=5
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
environment=NODE_ENV="production",DB_HOST="localhost",DB_PORT="5432",DB_USERNAME="admin",DB_PASSWORD="password123",DB_DATABASE="edumap_db",DB_SYNC="true",REDIS_HOST="localhost",REDIS_PORT="6379",JWT_SECRET="edumap_super_secret_jwt_key_2026",PORT="3001",MINIO_ACCESS_KEY="admin",MINIO_SECRET_KEY="password123",MINIO_ENDPOINT="localhost",MINIO_PORT="9000",AI_SERVICE_URL="http://127.0.0.1:8000"

[program:ai-service]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000
directory=/app/ai-service
priority=200
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/ai.log
stderr_logfile=/var/log/supervisor/ai.err

[program:frontend]
command=npm run start
directory=/app/frontend
priority=300
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/frontend.log
stderr_logfile=/var/log/supervisor/frontend.err

[program:nginx]
command=nginx -g "daemon off;"
priority=900
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err
